# utils/octopy_autostart/octopy.desktop
[Desktop Entry]
Type=Application
Name=Octopy
Exec=/home/{user}/Desktop/csomagolas.sh